# crontab
*/15 * * * * /openweather-influxdb-writer